# build/fxhub.service

[Unit]
Description=fxhub Server (10001)
After=default.target

[Service]
ExecStart=/opt/fxhub/fxhub 10001
Restart=always
RestartSec=1s

[Install]
WantedBy=default.target
